# root/etc/s6-overlay/s6-rc.d/svc-wireguard/finish
#!/bin/sh

info "Stopping Wireguard connection..."
wg show wg0
wg-quick down wg0
