# root/etc/s6-overlay/s6-rc.d/mtod/run
#!/bin/sh

cat << EOF
----------------------------------------------------
  ____            _                            _   
 |  _ \ _ __ ___ | |_ ___  _ __ _ __ ___ _ __ | |_ 
 | |_) | '__/ _ \| __/ _ \| '__| '__/ _ \ '_ \| __|
 |  __/| | | (_) | || (_) | |  | | |  __/ | | | |_ 
 |_|   |_|  \___/ \__\___/|_|  |_|  \___|_| |_|\__|

     Made with â¤ by https://github.com/bubuntux     
____________________________________________________   

ðŸ”§ Need help? https://github.com/bubuntux/protorrent/discussions/new
ðŸ› Bug? https://github.com/bubuntux/protorrent/issues/new
â˜• Discussion? https://github.com/bubuntux/protorrent/discussions/new
ðŸ’» Email? protorrent.limeade106@passmail.net
ðŸ’° Help me? https://github.com/sponsors/bubuntux

$(uname -a)
$(wg --version)
$(iptables --version)
$(qbittorrent-nox --version)
User UID: $(id -u abc)
User GID: $(id -g abc)

EOF
